# rootfs/usr/lib/systemd/system/install.service
[Unit]
Wants=init.scope

[Service]
Type=oneshot
ExecStart=/root/install
ExecStartPost=/root/cleanup
